# service/nginx/run
#!/bin/sh
mkdir -p /run/nginx/
sed -i s/ONION_BACKEND/$ONION_BACKEND/ /etc/service/nginx/nginx.conf
nginx -g "daemon off;" -c /etc/service/nginx/nginx.conf
